# common/php/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:php]
command=/usr/local/sbin/php-fpm -F
stopsignal=QUIT
stopasgroup=true
killasgroup=true
autorestart=true
stdout_logfile=/var/log/supervisor/php.out.log
stderr_logfile=/var/log/supervisor/php.err.log

[program:schedule]
directory=/app
command=/usr/local/bin/php artisan schedule:work --verbose --no-interaction
user=appuser
autostart=true
autorestart=true
stdout_logfile=/app/storage/logs/services/schedule.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5
stderr_logfile=/app/storage/logs/services/schedule-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
